# ansible/files/obnam.service
[Unit]
Description=Obnam server

[Service]
Type=simple
ConditionPathExists=/etc/obnam/server.yaml
ExecStart=/bin/obnam-server /etc/obnam/server.yaml

[Install]
WantedBy=multi-user.target
